# myapp.service
[Unit]
Description=Node.js Web Application
After=network.target mysql.service

[Service]
User=csye6225
Group=csye6225
WorkingDirectory=/home/csye6225/webapp
ExecStart=/usr/local/bin/nodemon /home/csye6225/webapp/app.js
Restart=always
StandardOutput=file:/var/log/myapp.log
StandardError=file:/var/log/myapp-error.log
EnvironmentFile=/home/csye6225/webapp/.env

[Install]
WantedBy=multi-user.target
