# myapp.service
[Unit]
Description=My Web Application
After=network.target

[Service]
User=csye6225
Group=csye6225
WorkingDirectory=/opt/app
EnvironmentFile=/opt/app/.env
ExecStart=/usr/bin/node app.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
